# docker/patches/rc.local
# 修复在某些docker环境下网络卡顿的问题, 如果非必要请勿删除!
cat > /etc/resolv.conf <<EOF
search lan
nameserver 127.0.0.1
options ndots:0
EOF

#
#
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

exit 0
